# deploy/ironmail-license.service
[Unit]
Description=IronMail License Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ironmail-license
EnvironmentFile=/etc/ironmail-license.env
ExecStart=/opt/ironmail-license/.venv/bin/uvicorn ironmail_license.app:app --host 127.0.0.1 --port 18081
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
